# systemd/backup.timer
[Unit]
Description=Daily backup timer for Linux Server Baseline

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
