# systemd/backup.service
[Unit]
Description=Linux Server Baseline Backup

[Service]
Type=oneshot
ExecStart=/opt/linux-server-baseline/scripts/backup.sh
